# ~/.config/autostart/clipboard-manager.desktop
[Desktop Entry]
Name=Clipboard Manager
Comment=Clipboard history manager for Linux
Exec=/usr/local/bin/clipboard-manager daemon
Icon=edit-copy
Terminal=false
Type=Application
Categories=Utility;System;
X-GNOME-Autostart-enabled=true
Hidden=false
NoDisplay=false
StartupNotify=false
X-GNOME-Autostart-Delay=3